# supervisord.conf
[unix_http_server]
file=/dev/shm/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
childlogdir=/var/log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///dev/shm/supervisor.sock

[program:ingoverkno]
command=/root/.rbenv/shims/bundle exec bin/dev
directory=/app
autostart=true
autorestart=false
stderr_logfile=/var/log/ingoverkno.err.log
stdout_logfile=/var/log/ingoverkno.out.log
environment=RAILS_ENV="development"